# packaging/etc/systemd/system/sshportal_clean_session_logs.timer
[Unit]
Description=Cleanup of SSHportal session logs
ConditionPathIsDirectory=/var/log/sshportal/session/

[Timer]
OnCalendar=weekly
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target